# etc/systemd/system/cigri.service.in
[Unit]
Description=CiGri computing grid server
After=postgresql.service

[Service]
ExecStartPre=+/usr/bin/install -o cigri -d %%PIDDIR%%
ExecStart=%%CIGRIDIR%%/modules/almighty.rb
PIDFile=%%PIDDIR%%/almighty.pid
User=%%CIGRIUSER%%
WorkingDirectory=/var/lib/cigri

[Install]
WantedBy=default.target
